# sky-pi.service
[Unit]
Description=sky-pi
After=syslog.target network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/robot
ExecStart=/home/pi/venv/bin/python /home/pi/robot/sky_pi.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
